# Scrapers/crontab
*/2 * * * * /app/run_scrapers.sh >> /tmp/cron.log 2>&1
